# autodeploy/stats.conf
# /etc/init/stats.conf

description "Meetin.gs statistics"
author      "Tuomas Starck <tuomas@meetin.gs>"

start on    (local-filesystems and net-device-up IFACE!=lo)
stop  on    shutdown

respawn
respawn limit 5 5

setuid nobody
setgid nogroup

script
  set -e
  FIFO=`mktemp --dry-run`
  mkfifo $FIFO
  (logger -t stats -p local5.info <$FIFO &)
  exec 2>$FIFO
  exec >$FIFO
  rm $FIFO

  . /etc/autodeploy.conf

  exec /usr/local/bin/stats
end script
